# modules/vagrant/master/resources/services/bootkube-wait.service
[Unit]
Before=bootkube.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /opt/tectonic/bootkube-wait.sh

[Install]
WantedBy=multi-user.target
