# modules/vagrant/master/resources/services/bootkube-wait.path
[Path]
PathExists=/opt/tectonic/bootkube.sh

[Install]
WantedBy=multi-user.target
